# various_config_files/scrapyd_upstart.conf
start on filesystem
setuid vagrant
setgid vagrant
script
        cd /vagrant/venv
        . bin/activate
        cd ..
        cd scrapyd
        scrapyd
end script